# ansible/roles/gogogadget_server/templates/gogogadget-server.service.j2
# Go Go Gadget Server Systemd Service File
[Unit]
Description=Go Go Gadget Server - Point Broadband Reachability Tool
After=network.target mariadb.service redis.service

[Service]
EnvironmentFile=/etc/gogogadget.env
ExecStart=/usr/local/bin/gogogadget --role server --port 8080 --config /etc/gogogadget/config.yml
Restart=always
User=root
RestartSec=5s

[Install]
WantedBy=multi-user.target